# supervisord.conf
[supervisord]
nodaemon=true

[program:ops-manager]
command=/opt/mongodb/mms/bin/mongodb-mms start
redirect_stderr=true
priority=1